# package/zyxel-streamboost/files_nbg6816/etc/hotplug.d/iface/99-qca-nss-ecm_ZyXEL
#!/bin/sh

if [ "$ACTION" = "ifup" ]; then
	/etc/init.d/qca-nss-ecm restart
fi
