# package/base-files/platform-dependency-files/NBG6817/bin/zyxel_switch_set
#!/bin/sh /etc/rc.common

enable(){
	setup_switch() { return 0; }

	include /lib/network
	setup_switch
}
